# deploy/systemd/llm-fineweb-stage-shard-loop.service
[Unit]
Description=LLM FineWeb Stage+Shard Loop Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-/etc/llm/llm.env
ExecStart=/bin/bash -lc 'cd __REPO_DIR__ && exec bash scripts/fineweb_stage_shard_loop.sh ${LLM_STAGE_SHARD_LOOP_ARGS:-"--hot-queue-min-files 12 --stage-max-files 10 --stage-copy-jobs 2 --stage-min-free-gib 80 --process-max-files 10 --shard-jobs 2 --auto-tune-shard-jobs --auto-tune-min-shard-jobs 1 --auto-tune-max-shard-jobs 4 --auto-tune-min-batch-seconds 300 --tokenizer-threads 10 --encode-batch-size 1024 --sleep-seconds 60 --shard-min-batch-size 512"}'
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=120
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
