# deploy/systemd/llm-train-supervisor.service
[Unit]
Description=LLM Train Supervisor (FineWeb 350BT)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-/etc/llm/llm.env
ExecStart=/bin/bash -lc 'cd __REPO_DIR__ && exec bash scripts/train_supervisor_rtx5070_350bt.sh ${LLM_TRAIN_SUPERVISOR_ARGS:-}'
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=120
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
