# deploy/systemd/llm-fineweb-prefetch.service
[Unit]
Description=LLM FineWeb Warm->Hot Prefetch Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-/etc/llm/llm.env
ExecStart=/bin/bash -lc 'cd __REPO_DIR__ && exec bash scripts/fineweb_prefetch_hot_queue.sh ${LLM_PREFETCH_ARGS:-}'
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=60
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
