# deploy/systemd/llm-shard-offload.service
[Unit]
Description=LLM Shard Bin Offload To Warm Storage
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-/etc/llm/llm.env
ExecStart=/bin/bash -lc 'cd __REPO_DIR__ && mkdir -p artifacts/reports && offload_args=${LLM_SHARD_OFFLOAD_ARGS:-"--shards-root data/shards_global/fineweb-global-bpe-v1 --warm-shards-root /mnt/ceph/llm/data/shards_global/fineweb-global-bpe-v1 --keep-local-batches 24 --target-free-gib 180 --max-batches 16 --disable-offloaded-manifests --require-trained-batches-file artifacts/reports/train_supervisor_phase1_talk/trained_batch_names.txt,artifacts/reports/train_supervisor_350bt/trained_batch_names.txt --skip-if-trained-file-missing --min-manifest-unique-input-files 510 --min-active-manifests 48 --min-active-train-tokens 40000000000"} && eval "set -- $offload_args" && exec flock -n artifacts/reports/shard_offload.lock env PYTHONPATH=src .venv/bin/python scripts/offload_shard_bins_to_warm.py "$@"'
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=true
