# deploy/systemd/llm-train-supervisor.service
[Unit]
Description=LLM Train Supervisor (FineWeb 350BT)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-/etc/llm/llm.env
ExecStart=/bin/bash -lc 'cd __REPO_DIR__ && exec bash scripts/train_supervisor_rtx5070_350bt.sh ${LLM_TRAIN_SUPERVISOR_ARGS:-"--step-chunk 2000 --poll-seconds 60 --batch-size 12 --target-effective-batch 24 --min-batch-size 6 --max-batch-size 20 --batch-step 2 --gpu-sample-seconds 1 --checkpoint-keep-last 6 --checkpoint-keep-every 10000 --ema-decay 0.999 --ema-update-every 1 --ema-start-step 12000 --eval-suite configs/eval/standard_prompt_suite_v2.json --no-train-fail-on-eval-regression"}'
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=120
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
