# deploy/systemd/llm-bad-parquet-revalidate.service
[Unit]
Description=LLM FineWeb Bad-Parquet Revalidate/Restage
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-/etc/llm/llm.env
ExecStart=/bin/bash -lc 'cd __REPO_DIR__ && PYTHONPATH=src .venv/bin/python scripts/revalidate_bad_parquet.py ${LLM_BAD_PARQUET_REVALIDATE_ARGS:-}'
NoNewPrivileges=true
