# deploy/systemd/user/llm-checkpoint-offload-prune.service
[Unit]
Description=LLM Checkpoint Offload And Local Prune
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-%h/.config/llm/llm.env
ExecStart=/bin/bash -lc 'cd __REPO_DIR__ && exec bash scripts/checkpoint_offload_prune.sh ${LLM_CHECKPOINT_OFFLOAD_ARGS:-}'
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=true
