# deploy/systemd/user/llm-hf-download-watchdog.service
[Unit]
Description=LLM Hugging Face Download Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-%h/.config/llm/llm.env
ExecStart=/bin/bash -lc 'cd __REPO_DIR__ && if [[ -z "${HF_TOKEN:-}" ]]; then echo "HF_TOKEN missing; set in ~/.config/llm/llm.env"; exit 2; fi; exec bash scripts/hf_download_watchdog.sh ${LLM_HF_WATCHDOG_ARGS:-}'
Restart=always
RestartSec=30
KillSignal=SIGINT
TimeoutStopSec=120
NoNewPrivileges=true

[Install]
WantedBy=default.target
